# misc/upstart/rpi-garage-door.conf
# node-site
#
# This upstart job configuration is managed by Puppet.
# Do not edit this file. Your changes will be lost.

description "Raspberry Pi Garage Door Controller"
author      "simianhacker"
version     "0.0.1"

start on runlevel [2345]
stop on runlevel [016]


respawn
respawn limit 5 10

console log

setuid pi
setgid gpio

chdir /home/pi/rpi-garage-door
exec /usr/local/bin/node index.js
